# files/supervisord.d/phantomjs.conf
[program:phantomjs]
command=phantomjs --webdriver=4444

