# files/supervisord.d/sshd.conf
[program:sshd]

# options: do not detach, errs to stderr, 
# logs to a log file (so syslog not needed)
# -e errs to stderr
command=/usr/sbin/sshd -D -E /var/log/sshd.log

